# templates/diamond.service.erb
[Unit]
Description=Diamond daemon
After=network.target

[Service]
ExecStart=<%= @diamond_path %> --foreground

[Install]
WantedBy=multi-user.target
