# templates/diamond.service.erb
[Unit]
Description=Diamond daemon
After=network.target

[Service]
ExecStart=<%= @diamond_path %> --foreground --log-stdout

[Install]
WantedBy=multi-user.target
